# payment.service
[Unit]
Description=Payment Service

[Service]
<<<<<<< HEAD
User=root
WorkingDirectory=/app
// highlight-start
Environment=CART_HOST=cart.shalinidaws88s.online
Environment=CART_PORT=8080
Environment=USER_HOST=user.shalinidaws88s.online
Environment=USER_PORT=8080
Environment=AMQP_HOST=<rabbitmq.shalinidaws88s.online
// highlight-end
Environment=AMQP_USER=roboshop
Environment=AMQP_PASS=roboshop123

ExecStart=/usr/local/bin/uwsgi --ini payment.ini
ExecStop=/bin/kill -9 $MAINPID
SyslogIdentifier=payment

[Install]
WantedBy=multi-user.target
=======
User=roboshop
Environment=REDIS_HOST=redis
Environment=MONGO_HOST=mongodb
Environment=RABBITMQ_HOST=rabbitmq
Environment=RABBITMQ_USER=roboshop
Environment=RABBITMQ_PASSWORD=roboshop
ExecStart=/usr/bin/python3 /app/payment.py
Restart=always

[Install]
WantedBy=multi-user.target
>>>>>>> 276ba82498b2fa04bd858c85a94b6abdbf55bd12
